# puppet/modules/yvonnesthing/templates/service.erb
#!/bin/bash
# chkconfig: 23456 20 80
# description: yvonnesthing

PATH=$PATH:/usr/local/bin

SERVERPATH=/usr/local/lib/node_modules/yvonnesthing
SERVERSCRIPT=$SERVERPATH/hello.js
RETVAL=0

start() { cd $SERVERPATH; forever start $SERVERSCRIPT; }
#cd /usr/local/lib/node_modules/yvonnesthing
#forever start /usr/local/lib/node_modules/yvonnesthing/hello.js
stop() { forever stop $SERVERSCRIPT; }
restart() { stop; start; }
status() {
        forever list | grep -i yvonnesthing
        RETVAL=$?
        if [ $RETVAL == 0 ]
        then
                echo "yvonnesthing is running"
        else
                echo "yvonnesthing is not running"
        fi
}

case "$1" in
        start)
                start
                ;;
        stop)
                stop
                ;;
        restart)
                restart
                ;;
        status)
                status
                ;;
        *)
                echo $"Usage: $0 {start|stop|restart|status}"
                RETVAL=2
esac

exit $RETVAL
